# xjtlu-timetable.service
[Unit]
Description=XJTLU Timetable (PocketBase)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/xjtlu-timetable/backend
ExecStart=/root/xjtlu-timetable/backend/pocketbase serve \
  --http=172.17.0.1:8091 \
  --publicDir=/root/xjtlu-timetable/frontend/dist \
  --dir=/root/xjtlu-timetable/backend/pb_data \
  --hooksDir=/root/xjtlu-timetable/backend/pb_hooks
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
